# etc/systemd/system/gamepads@.service
[Unit]
Description=Gamepads Service

[Service]
Type=simple
ExecStart=/home/pi/gamepads.sh gamepadConnect %i
